# app-i18n/rime-server/files/rime-server.conf
start on started system-services
stop on stopping system-services
respawn

exec minijail0 -u chronos -g chronos -- /bin/bash -c "cd /home/chronos/user/MyFiles/.rime && /usr/bin/rime_server > rime_server.log 2>&1"