# app-i18n/rime-server/files/rime-server.conf
start on start-user-session
stop on stopping ui
respawn

oom score -100

env RIME_PATH="/home/chronos/user/rime-conf/"

pre-start script
    [ ! -d ${RIME_PATH} ] && mkdir -p ${RIME_PATH} && cp -rp /usr/share/chromeos-assets/input_methods/rime-conf/* ${RIME_PATH} && sed -i s/12777420-82d5-45e1-830a-055dab6b307a/`/usr/bin/uuidgen`/g  ${RIME_PATH}/installation.yaml && touch ${RIME_PATH}/start.log && /bin/date >> ${RIME_PATH}/start.log && echo "install" >> ${RIME_PATH}/start.log
    diff -q /usr/share/chromeos-assets/input_methods/rime-conf /home/chronos/user/rime-conf | grep -v "rime-conf/build" | grep -v "Only in /home/chronos/user/rime-conf" | grep -v "Common subdirectories" | grep -v "user.yaml" && cp `ls -A /usr/share/chromeos-assets/input_methods/rime-conf/*.yaml | grep -v "user.yaml"` ${RIME_PATH} && /bin/date >> ${RIME_PATH}/start.log && echo "update" >> ${RIME_PATH}/start.log
    chown -R chronos:chronos ${RIME_PATH}
end script

post-stop script
    /bin/date >> ${RIME_PATH}/start.log
    cd ${RIME_PATH} && exec minijail0 -u chronos -g chronos -- /usr/bin/rime_dict_manager -s >> ${RIME_PATH}/start.log 2>&1
end script

exec minijail0 -u chronos -g chronos -- /usr/bin/rime_server ${RIME_PATH}
